# docs/launchd/io.goodkind.clyde.oauth-refresh.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>io.goodkind.clyde.oauth-refresh</string>
	<key>ProgramArguments</key>
	<array>
		<string>REPLACE_WITH_ABSOLUTE_PATH_TO_CLYDE</string>
		<string>oauth-refresh</string>
	</array>
	<key>StartInterval</key>
	<integer>1800</integer>
	<key>StandardOutPath</key>
	<string>REPLACE_WITH_HOME/.local/state/clyde/oauth-refresh.log</string>
	<key>StandardErrorPath</key>
	<string>REPLACE_WITH_HOME/.local/state/clyde/oauth-refresh.log</string>
</dict>
</plist>
