# packaging/systemd/clyde-daemon.service.in
[Unit]
Description=Clyde daemon
Documentation=https://github.com/agoodkind/clyde

[Service]
Type=simple
ExecStart=@@CLYDE_DAEMON_BIN@@ daemon
Restart=always
RestartSec=5
Environment=HOME=@@HOME@@
Environment=PATH=@@HOME@@/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
